# legacy/deploy/systemd/cowork-os-node.service
[Unit]
Description=CoWork OS (Node-Only, Headless)
After=network.target

[Service]
Type=simple

# Path where you cloned/built CoWork OS on the VPS.
WorkingDirectory=/opt/cowork-os

# Optional env file (LLM keys, data dir, control plane host/port, etc).
EnvironmentFile=-/etc/cowork-os.env

# Node-only daemon: no Electron/Xvfb required.
ExecStart=/usr/bin/node /opt/cowork-os/bin/coworkd-node.js

Restart=on-failure
RestartSec=5

# Run as an unprivileged user (recommended).
User=cowork
Group=cowork
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
CapabilityBoundingSet=
ReadWritePaths=/var/lib/cowork-os /srv/cowork/workspace

[Install]
WantedBy=multi-user.target
